# sonos_http_api/rootfs/etc/services.d/sonos-http-api/run
#!/usr/bin/with-contenv bashio
echo "Starting Sonos HTTP API..."
cd /app
exec node server.js
